# webapp/autoframe-web.service
[Unit]
Description=autoframe web UI
After=network.target

[Service]
Type=simple
User=pawkor
WorkingDirectory=/home/pawkor/scripts
ExecStart=/home/pawkor/highlight-env/bin/uvicorn webapp.server:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
KillMode=control-group
KillSignal=SIGKILL
TimeoutStopSec=3

[Install]
WantedBy=multi-user.target
